# sys-apps/pcsc-lite/files/_debian/postinst
#!/bin/sh

rc-update add pcscd default

