# sys-apps/timer_entropyd/files/_debian/postinst
#!/bin/sh

/sbin/rc-update add timer_entropyd default
