# sys-apps/timer_entropyd/files/_debian/prerm
#!/bin/sh

/sbin/rc-update del timer_entropyd default
